# files/92-dvita.rules
ACTION=="remove", SUBSYSTEM=="video4linux", RUN+="/home/pi/VitaDock/close.sh"
