# files/vita.service
[Service]
ExecStart=/bin/bash /home/pi/VitaDock/launch.sh
